# package-overrides/mirrorcell1m0-syslog/etc/config/crontab/root
#PATH=/sbin:/bin
#MAILTO=bburleson@lcogt.net
#MAILHOST=mail

0 * * * * /bin/ntpcheck
0-59 * * * * /bin/start_syslogd
